# pornote_supervisor.conf
[program:pornote]
command = gunicorn pornote:app
directory = /home/pi/napnac/pornote/
user = pi
